# packaging/systemd/boxup-backup-server@.timer
[Unit]
Description=Run the Boxup server backup for %i on its configured calendar

[Timer]
OnCalendar=*-*-* 04:00:00 UTC
RandomizedDelaySec=10min
Persistent=true
Unit=boxup-backup-now@%i.service

[Install]
WantedBy=timers.target
