# packaging/systemd/boxup-check@.service
[Unit]
Description=Check Boxup repository for %i
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/lib/boxup/boxup-root --config /etc/boxup/%i.toml check
Nice=10
IOSchedulingClass=idle
CPUWeight=20
IOWeight=20
CPUQuota=150%
MemoryHigh=2G
UMask=0077
NoNewPrivileges=yes
PrivateTmp=yes
PrivateDevices=yes
ProtectSystem=strict
ProtectHome=read-only
ProtectClock=yes
ProtectControlGroups=yes
ProtectHostname=yes
ProtectKernelLogs=yes
ProtectKernelModules=yes
ProtectKernelTunables=yes
ReadWritePaths=/var/lib/boxup/%i /var/cache/boxup/%i /var/lib/boxup-index/%i
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictNamespaces=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
LockPersonality=yes
SystemCallArchitectures=native
TimeoutStartSec=48h
TimeoutStopSec=10min
KillSignal=SIGTERM
KillMode=mixed
